# Genesis_Ecosystem/Enki-AI/linux/systemd/enki-hud.service
[Unit]
Description=Enki-AI Sovereign HUD (FastAPI / WebSocket)
Documentation=https://github.com/Cassai2026/Enki-AI
After=network.target

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment="PATH=__VENV_DIR__/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=__REPO_ROOT__"
EnvironmentFile=__REPO_ROOT__/.env
ExecStart=__VENV_DIR__/bin/python -m enki_ai.gui.hud_server --no-browser
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
